# debian/python-registryaggregator.service
[Unit]
Description=nmosregistration
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/bin/nmosregistration
